# dot_config/systemd/user/albert.service
[Unit]
Description=albert

[Service]
Environment=LANG=de_DE.utf8
ExecStart=/usr/bin/albert

[Install]
WantedBy=sway-session.target
