# dot_config/systemd/user/sway-session.target.wants/blueman-applet.service
[Unit]
Description=Network Manager Applet

[Service]
ExecStart=/usr/bin/blueman-applet
Restart=on-abort

[Install]
WantedBy=sway-session.target
