# dot_config/systemd/user/sway-session.target.wants/mega.service
[Unit]
Description=megasync

[Service]
ExecStart=/usr/bin/megasync

[Install]
WantedBy=sway-session.target
