# dot_config/systemd/user/sway-session.target.wants/nm-applet.service
[Unit]
Description=Network Manager Applet

[Service]
ExecStart=/usr/bin/nm-applet --indicator
Restart=on-abort

[Install]
WantedBy=sway-session.target
