# systemd/user/background.service
[Unit]
Description=set background image

[Service]
Type=oneshot
ExecStart=/usr/bin/feh --bg-fill --no-xinerama --randomize  -r ${HOME}/.wallpaper/
