# systemd/user/background.timer
[Unit]
Description=Change the wallpaper every 

[Timer]
OnCalendar=*:0/10

[Install]
WantedBy=x.target
