# systemd/user/compton.service
[Unit]
Description=compton

[Service]
ExecStart=/usr/bin/compton
Restart=on-abort

[Install]
WantedBy=x.target
