# systemd/user/dropbox.service
[Unit]
Description=Dropbox Network Storage Daemon

[Service]
ExecStart=/usr/bin/dropbox

[Install]
WantedBy=x.target
