# systemd/user/gnome-keyring.service
[Unit]
Description=gnome-keyring

[Service]
ExecStart=/usr/bin/gnome-keyring-daemon -f -s -c secrets
Restart=on-abort

[Install]
WantedBy=x.target
