# systemd/user/nm-applet.service
[Unit]
Description=Network Manager Applet

[Service]
ExecStart=/usr/bin/nm-applet
Restart=on-abort

[Install]
WantedBy=x.target
