# systemd/user/redshift.service
[Unit]
Description=Start Redshift Applet

[Service]
ExecStart=/usr/bin/redshift-gtk
Restart=on-abort

[Install]
WantedBy=x.target
