# systemd/user/x.target
[Unit]
Description=target for user services requiere X
Requires=default.target
After=default.target
